# systemd/.config/systemd/user/aria2cd.service
[Unit]
Description=Aria2 Service
After=network-online.target

[Service]
Type=forking
ExecStart=/usr/bin/aria2c -D

[Install]
WantedBy=default.target
