# systemd/.config/systemd/user/fava.service
[Unit]
Description=Fava daemon
After=network.target

[Service]
ExecStart=/usr/bin/fava %h/accounting/ToC.bean
Restart=on-failure

[Install]
WantedBy=default.target
