# systemd/.config/systemd/user/network-online.target
[Unit]
Description=Network is Online

[Install]
WantedBy=default.target
