# systemd/.config/systemd/user/refresh-gmail-xoauth-token.timer
[Unit]
Description=Refresh GMail XOAuth token daily

[Timer]
OnCalendar=*:0/15
Persistent=true

[Install]
WantedBy=timers.target
